# package-templates/dropbear/etc/profile.d/dropbear_git.sh
#!/bin/dash

# /etc/profile.d/dropbear_git.sh: a script for setting the Git's SSH client

export GIT_SSH=/usr/bin/dropbear_ssh