# package-templates/shahor-desktop/etc/profile.d/rxvt-unicode.sh
#!/bin/sh

# /etc/profile.d/rxvt-unicode.sh: environmental variables related to
#                                 rxvt-unicode

# set the urxvtd socket path
export RXVT_SOCKET="$HOME/.urxvt_socket"