# skeleton/rootfs/var/spool/cron/crontabs/root
# update locate's database, every day at midnight
0 0 * * * updatedb
